# systemd/flamos.service
[Unit]
Description=FlaMoS Daemon
After=multi-user.target

[Service]
Type=idle
WorkingDirectory=/home/flamos/FlaMoS
ExecStart=/home/flamos/venv/bin/python3 flamosd.py

[Install]
WantedBy=multi-user.target